# crontab
* * * * * java -cp application.jar io.github.raniagus.example.bootstrap.Bootstrap
